# examples/s6/mdevd/run
#!/bin/execlineb -P
exec -c
mdevd-netlink -d 3
